# 3.0/config/etc/supervisor/conf.d/mongodb.conf
[program:mongodb]
command=/usr/bin/mongod --smallfiles --auth
redirect_stderr=true
stopsignal=TERM
autorestart=true
